# k8s-sample/crontab
* * * * * cd /sechecker/ && ./sechecker > /dev/stdout 2>&1